# backend/crontab
# MPZ 공공데이터 동기화 스케줄 (UTC 기준, KST = UTC+9)

# 매일 03:00 KST (18:00 UTC) — 증분 동기화 (어제~오늘)
0 18 * * * cd /backend && /opt/venv/bin/python manage.py sync_public_data --strategy incremental --days 2

# 매주 일요일 04:00 KST (19:00 UTC) — 상태만 동기화 (이미지 스킵)
0 19 * * 0 cd /backend && /opt/venv/bin/python manage.py sync_public_data --strategy status_sync

# 매월 1일 05:00 KST (20:00 UTC) — 전체 동기화 (최근 90일)
0 20 1 * * cd /backend && /opt/venv/bin/python manage.py sync_public_data --strategy full --days 90
